# root/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

USER_NAME=${USER_NAME:-linuxserver.io}

exec \
        s6-setuidgid ${USER_NAME} /usr/sbin/sshd -D -e -p 2222 -E /config/logs/openssh/openssh.log
