# root/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

exec \
        /usr/sbin/sshd -D -e