# root/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

USER_NAME=${USER_NAME:-linuxserver.io}

exec 2>&1 \
        s6-setuidgid ${USER_NAME} /usr/sbin/sshd -D -e -p 2222
