# root/etc/s6-overlay/s6-rc.d/init-openssh-server-config/up
/etc/s6-overlay/s6-rc.d/init-openssh-server-config/run
