# root/etc/s6-overlay/s6-rc.d/init-nuj-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "running nuj stuff"
aws s3 sync s3://"${CONFIG_BUCKET}" /config
echo "nuj config applied"
