# root/etc/s6-overlay/s6-rc.d/svc-openssh-server/run
#!/usr/bin/with-contenv bash

USER_NAME=${USER_NAME:-linuxserver.io}
OPENSSH_PORT=${OPENSSH_PORT:-2222}

exec 2>&1 \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${OPENSSH_PORT}" \
    s6-setuidgid "${USER_NAME}" /usr/sbin/sshd.pam -D -e -p ${OPENSSH_PORT}
